# .config/systemd/user/battery-warn.timer

[Unit]
Description=Battery Warning Timer

[Timer]
OnBootSec=60
OnUnitActiveSec=2min

[Install]
WantedBy=timers.target
