# .config/systemd/user/battery-warn.service
[Unit]
Description=Battery Warning

[Service]
Type=oneshot
ExecStart=/home/jjoyner/.local/bin/battery-warn.sh
